# dist/openeis-support/systemd/openeis.socket
[Unit]
Description=OpenEIS uWSGI Server

[Socket]
ListenStream=/var/run/openeis/uwsgi.sock

[Install]
WantedBy=sockets.target
